# roles/live-srv/files/bar-backup.service
[Unit]
Description=PostgreSQL backup for BAR Live Services

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-take.sh
