# roles/live-srv/files/bar-backup.timer
[Unit]
Description=Daily PostgreSQL backup for BAR Live Services

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
